# deploy/college-ai-api.service
[Unit]
Description=College AI FastAPI Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ec2-user/college-ai
EnvironmentFile=/home/ec2-user/college-ai/.env
ExecStart=/home/ec2-user/college-ai/venv/bin/python -m uvicorn college_ai.api.app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
